# herramientas/reporte/report.conf
[program:report]
command=/usr/bin/wine /app/report.exe
priority=100
stdout_logfile=/var/log/report.stdout.log
stderr_logfile=/var/log/report.stderr.log